# my_crontab.crontab

* 9 * * * cd /Users/sunyunxian/workspace/code/dotfiles; git status; git fetch
